# roles/jontheniceguy.rpirouter/files/libComposite.service
[Unit]
Description=Start libComposite

[Service]
Type=oneshot
ExecStart=/opt/lib_composite.sh

[Install]
WantedBy=multi-user.target